# maintenance/system_cleanup.service
[Unit]
Description=Piphawk AI system cleanup

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/piphawk/maintenance/system_cleanup.py --yes

[Install]
WantedBy=multi-user.target
